# etc/cron.d/wordpress
5 *	* * *	fcgi-php	cd /wordpress && php ./wp-cron.php
